# opencode-orchestrator.service
[Unit]
Description=Opencode Stateless Orchestrator
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/opencode-orchestration
# Load environment variables from .env file
EnvironmentFile=/root/opencode-orchestration/.env
# Run the bot using the shell script to ensure venv is used
ExecStart=/root/opencode-orchestration/run_bot.sh
# Restart automatically if it crashes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
